# scripts/com.stock.dailyupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
E1: macOS launchd configuration for scheduled stock predictor updates

INSTALLATION:
1. Copy this file to ~/Library/LaunchAgents/
   cp scripts/com.stock.dailyupdate.plist ~/Library/LaunchAgents/

2. Modify the paths below to match your setup:
   - WorkingDirectory: Your project path
   - ProgramArguments: Path to daily_update.sh

3. Load the agent:
   launchctl load ~/Library/LaunchAgents/com.stock.dailyupdate.plist

4. To unload:
   launchctl unload ~/Library/LaunchAgents/com.stock.dailyupdate.plist

5. Check status:
   launchctl list | grep stock

SCHEDULE: Runs at 6:30 AM every weekday (Mon-Fri)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stock.dailyupdate</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yadu/Documents/doc/stock/scripts/daily_update.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/yadu/Documents/doc/stock</string>
    
    <!-- Run at 6:30 AM on weekdays -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- Tuesday -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- Wednesday -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- Thursday -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- Friday -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHON_PATH</key>
        <string>/usr/bin/python3</string>
    </dict>
    
    <!-- Standard output and error logs -->
    <key>StandardOutPath</key>
    <string>/Users/yadu/Documents/doc/stock/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/yadu/Documents/doc/stock/logs/launchd_stderr.log</string>
    
    <!-- Run missed jobs if computer was asleep -->
    <key>StartInterval</key>
    <integer>0</integer>
</dict>
</plist>
